# conf/.config/systemd/user/librefang.service
[Unit]
Description=LibreFang Agent OS Daemon
Documentation=https://librefang.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=${HOME}/.bin/librefang start --foreground
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
WorkingDirectory=${HOME}/.librefang
EnvironmentFile=-${HOME}/.librefang/env
EnvironmentFile=-${HOME}/.librefang/secrets.env

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=default.target
