# etc/supervisor.d/fcitx5.ini
[program:fcitx5]
command = fcitx5 -D --replace
environment=HOME=/home/imxieke,XDG_RUNTIME_DIR="/run/user/1001"
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = imxieke
stdout_logfile = /var/log/fcitx5.out.log
stderr_logfile = /var/log/fcitx5.err.log